# practice/n_pre_school/beginner/chapter2/crontab
0 * * * * /work/bot/niconico-ranking.sh
